# pi-installer-backend.service
# PI-Installer Backend (nur API auf :8000)
# Für DSI Radio, API-Zugriff etc. Unabhängig vom Frontend.
# Installieren: siehe scripts/install-backend-service.sh oder create_installer.sh
# Nach Backend-Update: sudo systemctl restart pi-installer-backend

[Unit]
Description=PI-Installer Backend (API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/scripts/start-backend.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PI_INSTALLER_DEV=0
# Härtung (bestehende Optionen beibehalten + erweitert)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
MemoryMax=512M
LimitNOFILE=4096
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
ReadWritePaths={{INSTALL_DIR}} /var/lib/pi-installer /tmp

[Install]
WantedBy=multi-user.target
