# setuphelfer-backend.service
# Setuphelfer Backend (nur API auf :8000)
# Installieren: siehe scripts/install-backend-service.sh oder create_installer.sh

[Unit]
Description=Setuphelfer Backend (API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/scripts/start-backend.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PI_INSTALLER_DEV=0
Environment=APP_EDITION=release
Environment=PIP_CACHE_DIR={{INSTALL_DIR}}/.pip-cache
Environment="SETUPHELFER_DIR={{INSTALL_DIR}}"
Environment="PI_INSTALLER_DIR={{INSTALL_DIR}}"
Environment="SETUPHELFER_CONFIG_DIR={{PI_INSTALLER_CONFIG_DIR}}"
Environment="PI_INSTALLER_CONFIG_DIR={{PI_INSTALLER_CONFIG_DIR}}"
Environment="SETUPHELFER_LOG_DIR={{PI_INSTALLER_LOG_DIR}}"
Environment="PI_INSTALLER_LOG_DIR={{PI_INSTALLER_LOG_DIR}}"
Environment="SETUPHELFER_STATE_DIR={{PI_INSTALLER_STATE_DIR}}"
Environment="PI_INSTALLER_STATE_DIR={{PI_INSTALLER_STATE_DIR}}"
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
MemoryMax=512M
LimitNOFILE=4096
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
ReadWritePaths={{INSTALL_DIR}} /etc/setuphelfer /var/log/setuphelfer /var/lib/setuphelfer /tmp

[Install]
WantedBy=multi-user.target
