# packaging/systemd/setuphelfer-backend.service.d/backup-target-br001.conf
# BR-001: externes Ziel unter /media/setuphelfer/… (ProtectSystem=strict).
# Installieren nach /etc/systemd/system/setuphelfer-backend.service.d/ und daemon-reload + restart.
[Service]
ReadWritePaths=/media/setuphelfer
