# packaging/systemd/setuphelfer-backend.service.d/backup-target-br001.conf
[Service]
ReadWritePaths=/media/setuphelfer/br001
