# debian/postinst
#!/bin/sh
# set -e entfernt: Fehlerbehandlung explizit, damit dpkg nicht bei jedem Fehler abbricht

INSTALL_DIR="/opt/setuphelfer"
CONFIG_DIR="/etc/setuphelfer"
LOG_DIR="/var/log/setuphelfer"
STATE_DIR="/var/lib/setuphelfer"
SERVICE_USER="setuphelfer"

case "$1" in
  configure)
    # Alte systemd-Units stilllegen (Migration von pi-installer → setuphelfer)
    for old in pi-installer.service pi-installer-backend.service; do
      systemctl stop "$old" 2>/dev/null || true
      systemctl disable "$old" 2>/dev/null || true
    done

    # Service-User und Gruppe anlegen
    if ! getent passwd "$SERVICE_USER" >/dev/null 2>&1; then
      if ! getent group "$SERVICE_USER" >/dev/null 2>&1; then
        groupadd --system "$SERVICE_USER" 2>/dev/null || true
      fi
      if ! adduser --system --no-create-home --disabled-login --gecos "Setuphelfer Service" --ingroup "$SERVICE_USER" "$SERVICE_USER" 2>/dev/null; then
        useradd --system --no-create-home --comment "Setuphelfer Service" --gid "$SERVICE_USER" "$SERVICE_USER" 2>/dev/null || \
        useradd --system --no-create-home --comment "Setuphelfer Service" "$SERVICE_USER" 2>/dev/null || true
      fi
    fi

    SERVICE_GROUP=$(id -gn "$SERVICE_USER" 2>/dev/null || echo "$SERVICE_USER")

    for _svc in /etc/systemd/system/setuphelfer.service /etc/systemd/system/setuphelfer-backend.service; do
      if [ -f "$_svc" ]; then
        sed -i "s/^Group=.*/Group=$SERVICE_GROUP/" "$_svc" 2>/dev/null || true
      fi
    done

    mkdir -p "$CONFIG_DIR" "$LOG_DIR" "$STATE_DIR"

    # Konfiguration von Legacy übernehmen (einmalig), wenn Ziel noch leer
    if [ ! -f "$CONFIG_DIR/config.json" ] && [ -f "/etc/pi-installer/config.json" ]; then
      cp -a "/etc/pi-installer/config.json" "$CONFIG_DIR/" 2>/dev/null || true
    fi
    if [ ! -f "$CONFIG_DIR/backup.json" ] && [ -f "/etc/pi-installer/backup.json" ]; then
      cp -a "/etc/pi-installer/backup.json" "$CONFIG_DIR/" 2>/dev/null || true
    fi

    # Temporär auf root setzen, damit root schreiben kann
    [ -d "$INSTALL_DIR" ] && chown -R root:root "$INSTALL_DIR" 2>/dev/null || true

    if [ ! -f "$CONFIG_DIR/config.yaml" ]; then
      cat > "$CONFIG_DIR/config.yaml" << 'EOF'
install_dir: /opt/setuphelfer
config_dir: /etc/setuphelfer
log_dir: /var/log/setuphelfer
state_dir: /var/lib/setuphelfer
backend:
  host: 0.0.0.0
  port: 8000
frontend:
  port: 3001
EOF
    fi

    if [ -d "$INSTALL_DIR/backend" ] && [ ! -d "$INSTALL_DIR/backend/venv" ]; then
      PYTHON=""
      for py in python3.12 python3.11 python3.10 python3.9 python3; do
        if command -v "$py" >/dev/null 2>&1; then
          if "$py" -c "import sys; exit(0 if sys.version_info >= (3, 9) else 1)" 2>/dev/null; then
            PYTHON="$py"
            break
          fi
        fi
      done
      if [ -n "$PYTHON" ]; then
        mkdir -p "$INSTALL_DIR/.pip-cache"
        cd "$INSTALL_DIR/backend" 2>/dev/null || true
        if [ -f "requirements.txt" ]; then
          if "$PYTHON" -m venv venv 2>/dev/null && \
             export PIP_CACHE_DIR="$INSTALL_DIR/.pip-cache" && \
             ./venv/bin/pip install --upgrade pip -q 2>/dev/null && \
             ./venv/bin/pip install -q -r requirements.txt 2>&1 | grep -v "already satisfied" >/dev/null 2>&1; then
            chown -R "$SERVICE_USER:$SERVICE_GROUP" "$INSTALL_DIR/backend" "$INSTALL_DIR/.pip-cache" 2>/dev/null || true
          else
            echo "Warnung: Backend-Venv konnte nicht erstellt werden. Später: cd $INSTALL_DIR/backend && python3 -m venv venv && ./venv/bin/pip install -r requirements.txt" >&2
          fi
        fi
      fi
    fi

    if [ -d "$INSTALL_DIR/frontend" ] && [ ! -d "$INSTALL_DIR/frontend/node_modules" ]; then
      if command -v npm >/dev/null 2>&1; then
        cd "$INSTALL_DIR/frontend" 2>/dev/null || true
        if [ -f "package.json" ]; then
          export npm_config_cache="$INSTALL_DIR/.npm-cache"
          mkdir -p "$npm_config_cache"
          if npm install --silent 2>&1 | grep -v "npm WARN" >/dev/null 2>&1; then
            chown -R "$SERVICE_USER:$SERVICE_GROUP" "$INSTALL_DIR/frontend" "$INSTALL_DIR/.npm-cache" 2>/dev/null || true
          else
            echo "Warnung: Frontend npm install konnte nicht ausgeführt werden. Später: cd $INSTALL_DIR/frontend && npm install" >&2
          fi
        fi
      fi
    fi

    [ -d "$INSTALL_DIR" ] && chown -R "$SERVICE_USER:$SERVICE_GROUP" "$INSTALL_DIR" 2>/dev/null || true
    [ -d "$CONFIG_DIR" ] && chown -R "$SERVICE_USER:$SERVICE_GROUP" "$CONFIG_DIR" 2>/dev/null || true
    [ -d "$LOG_DIR" ] && chown -R "$SERVICE_USER:$SERVICE_GROUP" "$LOG_DIR" 2>/dev/null || true
    [ -d "$STATE_DIR" ] && chown -R "$SERVICE_USER:$SERVICE_GROUP" "$STATE_DIR" 2>/dev/null || true
    [ -d "$INSTALL_DIR" ] && chmod 755 "$INSTALL_DIR" 2>/dev/null || true
    [ -d "$CONFIG_DIR" ] && chmod 755 "$CONFIG_DIR" 2>/dev/null || true
    [ -d "$LOG_DIR" ] && chmod 755 "$LOG_DIR" 2>/dev/null || true
    [ -d "$STATE_DIR" ] && chmod 755 "$STATE_DIR" 2>/dev/null || true

    if [ ! -x "$INSTALL_DIR/frontend/src-tauri/target/release/pi-installer" ]; then
      BUILD_USER="${SUDO_USER:-}"
      if [ -n "$BUILD_USER" ] && su -s /bin/bash "$BUILD_USER" -c "command -v cargo" >/dev/null 2>&1; then
        echo "Tauri-Binary fehlt. Versuche es zu bauen (kann einige Minuten dauern)..."
        chown -R "$BUILD_USER:$BUILD_USER" "$INSTALL_DIR/frontend" 2>/dev/null || true
        if su -s /bin/bash "$BUILD_USER" -c "cd $INSTALL_DIR/frontend && export GDK_BACKEND=x11 && npm run tauri:build" 2>/dev/null; then
          echo "Tauri-Binary erstellt – App-Fenster verfügbar."
        else
          echo "⚠️  Tauri-Build fehlgeschlagen. App-Fenster nicht verfügbar (Browser-Modus funktioniert weiterhin)."
        fi
        chown -R "$SERVICE_USER:$SERVICE_GROUP" "$INSTALL_DIR/frontend" 2>/dev/null || true
      else
        echo "ℹ️  Tauri-Binary nicht im Paket und kein Cargo verfügbar. Browser-Modus wird verwendet."
      fi
    else
      echo "✅ Tauri-Binary gefunden – App-Fenster verfügbar."
    fi

    systemctl daemon-reload 2>/dev/null || true
    systemctl enable setuphelfer-backend.service 2>/dev/null || true
    systemctl enable setuphelfer.service 2>/dev/null || true
    systemctl start setuphelfer-backend.service 2>/dev/null || true
    systemctl start setuphelfer.service 2>/dev/null || true

    if command -v update-desktop-database >/dev/null 2>&1; then
      update-desktop-database -q /usr/share/applications 2>/dev/null || true
    fi

    echo "Setuphelfer wurde installiert. Startmenü: SetupHelfer / SetupHelfer (Browser)."
    echo "Schreibtisch-Starter: bash $INSTALL_DIR/scripts/desktop-setuphelfer-launcher-anlegen.sh"
    echo "Services: systemctl status setuphelfer-backend setuphelfer"
    ;;
  abort-upgrade|abort-remove|abort-deconfigure)
    ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

exit 0
