# debian/setuphelfer.service
[Unit]
Description=Setuphelfer Web-UI (vite preview; Backend in setuphelfer-backend.service)
Documentation=https://github.com/Yammato-2035/piinstaller
After=network-online.target setuphelfer-backend.service
Wants=network-online.target
Requires=setuphelfer-backend.service

[Service]
Type=simple
User=setuphelfer
# Group wird in postinst auf die tatsächliche Gruppe gesetzt (kann nogroup sein)
Group=nogroup
WorkingDirectory=/opt/setuphelfer
ExecStart=/opt/setuphelfer/scripts/start-browser-production.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="SETUPHELFER_DIR=/opt/setuphelfer"
Environment="PI_INSTALLER_DIR=/opt/setuphelfer"
Environment="SETUPHELFER_CONFIG_DIR=/etc/setuphelfer"
Environment="PI_INSTALLER_CONFIG_DIR=/etc/setuphelfer"
Environment="SETUPHELFER_LOG_DIR=/var/log/setuphelfer"
Environment="PI_INSTALLER_LOG_DIR=/var/log/setuphelfer"
Environment="PI_INSTALLER_DEV=0"
Environment="APP_EDITION=release"
Environment="PIP_CACHE_DIR=/opt/setuphelfer/.pip-cache"
Environment="npm_config_cache=/opt/setuphelfer/.npm-cache"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/setuphelfer /etc/setuphelfer /var/log/setuphelfer /var/lib/setuphelfer /tmp
MemoryMax=1024M
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
