# scripts/rescue-live/image/systemd/setuphelfer-rescue-ui.service
[Unit]
Description=Setuphelfer Rescue React Shell (offline-first)
After=local-fs.target setuphelfer-rescue-state.service
ConditionPathExists=/usr/share/setuphelfer/rescue/ui/rescue.html
ConditionKernelCommandLine=!setuphelfer_start_assistant=1

[Service]
Type=simple
Environment=TERM=linux
ExecStart=/usr/local/sbin/setuphelfer-rescue-ui-launch
Restart=on-failure
RestartSec=3
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
