# scripts/rescue-live/image/systemd/setuphelfer-rescue-ui.service
[Unit]
Description=Setuphelfer Rescue React Shell (offline-first)
After=local-fs.target setuphelfer-rescue-state.service
ConditionPathExists=/usr/share/setuphelfer/rescue/ui/rescue.html
# The graphical Chromium kiosk needs a running X server. On real hardware (and
# under nomodeset) there is no X here, so Chromium exits with "Missing X server"
# and Restart= turns it into a crash-loop that spams tty1 and makes the text TUI
# unusable. The kiosk is therefore OPT-IN via setuphelfer_kiosk=1; by default the
# text start-assistant TUI owns tty1.
ConditionKernelCommandLine=setuphelfer_kiosk=1

[Service]
Type=simple
Environment=TERM=linux
ExecStart=/usr/local/sbin/setuphelfer-rescue-ui-launch
Restart=on-failure
RestartSec=3
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
