# scripts/rescue-live/image/systemd/setuphelfer-rescue-ui.service
[Unit]
Description=Setuphelfer Rescue React Shell (offline-first)
After=local-fs.target setuphelfer-rescue-state.service setuphelfer-backend.service
Wants=setuphelfer-backend.service
ConditionPathExists=/usr/share/setuphelfer/rescue/ui/rescue.html

[Service]
Type=simple
Environment=TERM=linux
ExecStartPre=/usr/local/sbin/setuphelfer-rescue-backend-start.sh
ExecStart=/usr/local/sbin/setuphelfer-rescue-gui-start
Restart=on-failure
RestartSec=5
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
