# packaging/systemd/setuphelfer-backup@.service.d/backup-target-br001.conf
# BR-001: externes Ziel unter /media/setuphelfer/… (Betreiber, ohne /media/<login>/…-Abhängigkeit).
# systemd führt ReadWritePaths aus Unit + Drop-ins zusammen.
[Service]
ReadWritePaths=/media/setuphelfer
