# packaging/systemd/setuphelfer-backup@.service.d/backup-target-br001.conf
# BR-001: externes Ziel setuphelfer-back1 (/dev/sda1) — Betreiberfreigabe
[Service]
ReadWritePaths=/media/gabriel/setuphelfer-back /media/gabriel/setuphelfer-back1
Environment=PYTHONPATH=/opt/setuphelfer/backend
