# packaging/systemd/setuphelfer-backup@.service.d/notification-env.conf.example
# Optional drop-in: same notification.env as setuphelfer-backend.service
# Install: sudo cp notification-env.conf.example /etc/systemd/system/setuphelfer-backup@.service.d/notification-env.conf
#          sudo systemctl daemon-reload
[Service]
EnvironmentFile=-/etc/setuphelfer/notification.env
