# build/rescue/profiles/developer-qemu/systemd/setuphelfer-dev-agent.service
[Unit]
Description=Setuphelfer Development Agent (Rescue Developer QEMU Lab)
After=network-online.target setuphelfer-backend.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/setuphelfer/setuphelfer-dev-agent.env
WorkingDirectory=/opt/setuphelfer-rescue/backend
Environment=PYTHONPATH=/opt/setuphelfer-rescue/backend:/opt/setuphelfer-rescue
Environment=SETUPHELFER_RESCUE_ROOT=/opt/setuphelfer-rescue
ExecStart=/usr/bin/python3 -m devserver_agent.cli --send --json --qemu-host-fallback
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/setuphelfer-rescue/docs/evidence/runtime-results/dev-agent-spool /tmp /run

[Install]
WantedBy=multi-user.target
